# paste.bin.ingest.cmd.crontab
* */4 * * * dotnet paste.bin.ingest.cmd.dll
