# paste.bin.ingest.cmd.crontab
0 */4 * * * /bin/paste.bin.ingest/paste.bin.ingest.cmd
